# Arch Linux/home/ni/.config/autostart/Bing Photo of the Day.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=python3 "/home/ni/Documents/Repositories/Miscellaneous-Projects/Bing Photo of the Day/solution.py"
Name=Bing Photo of the Day
